# penstation.service
[Unit]
Description=PENSTATION — Autonomous Network Security Station
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/pi/penstation
Environment="PATH=/home/pi/penstation/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/pi/penstation/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=penstation

[Install]
WantedBy=multi-user.target
